# systemd/spirit-island.service
[Unit]
Description=Spirit Island web site

[Service]
WorkingDirectory=/home/ubuntu/spirit-island/
ExecStart=/home/ubuntu/spirit-island/run.sh

[Install]
WantedBy=default.target
